# conf/systemd/getonionbridges.timer
[Unit]
Description="Run script every day at 2pm"

[Timer]
Unit=getonionbridges.service
OnCalendar=*-*-* 2:00:00

[Install]
WantedBy=default.target
